# base/dotfiles/crontab
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# TODO: Is this still needed?
#
# Clean up /dev/shm every minute.
* * * * * $HOME/bin/clean-dev-shm

