# common-text/dotfiles/crontab
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Run logrotate against ~/log every hour.
0 * * * * logrotate --state $HOME/.logrotate.status $HOME/.logrotate.conf
