# src/dotfiles/crontab
# Frequently make sure that the right desktop wallpaper is shown.
* * * * * $HOME/bin/fehbg
